# scripts/audit-poste/Collecte-AuditPoste-Task.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Modele de tache planifiee — Collecte audit postes Parc.
  A adapter : chemin UNC du script (\\serveur\IT$\audit-poste\...).

  Import manuel (test sur 1 PC) :
    schtasks /Create /TN "Parc\CollecteAuditPoste" /XML "Collecte-AuditPoste-Task.xml" /F

  Ou via GPO (voir DEPLOY-GPO.md).
-->
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Date>2026-07-23T08:00:00</Date>
    <Author>Parc Informatique</Author>
    <Description>Collecte quotidienne des audits postes vers l API Parc (POST /api/audit).</Description>
    <URI>\Parc\CollecteAuditPoste</URI>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-07-23T08:15:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
      <RandomDelay>PT30M</RandomDelay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>
      <RunLevel>HighestAvailable</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT15M</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -NonInteractive -ExecutionPolicy Bypass -File "\\SERVEUR\IT$\audit-poste\Collecte-AuditPoste.ps1"</Arguments>
    </Exec>
  </Actions>
</Task>
